# files/warpd.service
[Install]
WantedBy=default.target

[Service]
Type=simple
# `-f` option is to make output go to journald
ExecStart=warpd -f
